# supervisord.conf
[program:node]
command = node server.js
directory = /home/vcap/app
